# rootfs/etc/services.d/minidsp/run
#!/usr/bin/with-contenv bash
source /opt/venv/bin/activate
exec python /app/api.py
