# rootfs/etc/services.d/minidsp/run
#!/command/with-contenv
exec /opt/venv/bin/python3 /app/api.py
