# rootfs/etc/services.d/minidsp/run
#!/command/with-contenv
exec python3 /app/api.py
