# rootfs/etc/services.d/minidsp/run
#!/command/with-contenv
/opt/venv/bin/python /app/api.py